# multi-hid-gadget.service
[Unit]
Description=Setup USB Gadget for Multi HID Proxy
DefaultDependencies=no
After=systemd-modules-load.service
Before=sysinit.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/setup_hid_gadget.sh

[Install]
WantedBy=sysinit.target

